# packages/jahia/migrations/v38_assets/tomcat_tmp_cleanup.service
[Unit]
Description=Clean tomcat folder

[Service]
Type=oneshot
ExecStart=/bin/bash -c "find /opt/tomcat/temp/* -not -path '*/.*' -type f -mtime +7 -exec bash -c 'fuser -s {} || rm {}' \; \
; find /opt/tomcat/temp/* -empty -not -path '*/.*' -type d -delete"
